# chroot-layer/etc/profile.d/k3s.sh
export K3S_DATASTORE_ENDPOINT="mysql://k3s:secret@tcp(nas:3306)/k3s"
export K3S_TOKEN="unsafe"
export INSTALL_K3S_VERSION="v1.24.3+k3s1"
