# k3s-agent/etc/systemd/system/etc-rancher-node.automount
[Unit]
Description=Persistent node settings automount
Before=k3s-agent.service

[Automount]
Where=/etc/rancher/node

[Install]
WantedBy=multi-user.target
